# ugent-line-proxy.service
[Unit]
Description=UGENT LINE Proxy Server
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/var/lib/ugent-line-proxy
EnvironmentFile=/etc/ugent-line-proxy/.env
ExecStart=/usr/bin/ugent-line-proxy
Restart=on-failure
RestartSec=5
TimeoutStopSec=30
LimitNOFILE=65536

# Security hardening (relaxed for root user)
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
